# files/systemd/system/sunshine.service
[Unit]
Description=Sunshine Game Streaming Service

StartLimitIntervalSec=500
StartLimitBurst=5

After=graphical-session.target

[Service]
ExecStart=/usr/bin/sunshine  '~/.config/sunshine/sunshine.conf'
ExecStartPre=/usr/bin/sleep 10
Restart=on-failure
RestartSec=5s


[Install]
WantedBy=graphical-session.target
